# deploy/aurora.service
[Unit]
Description=Aurora Linux EDR Agent
After=network.target

[Service]
Type=simple
Environment=AURORA_RULES_DIR=/opt/aurora-linux/sigma-rules/linux
Environment=AURORA_LOG_FILE=/var/log/aurora-linux/aurora.log
Environment=AURORA_OUTPUT_FLAG=--json
Environment=AURORA_EXTRA_ARGS=
EnvironmentFile=-/opt/aurora-linux/config/aurora.env
ExecStart=/bin/bash -c 'exec /opt/aurora-linux/aurora \
    --rules "$AURORA_RULES_DIR" \
    --logfile "$AURORA_LOG_FILE" \
    $AURORA_OUTPUT_FLAG \
    $AURORA_EXTRA_ARGS'
Restart=on-failure
RestartSec=5

# Privileges
User=root
# Alternative non-root setup:
# User=aurora
# AmbientCapabilities=CAP_BPF CAP_PERFMON CAP_SYS_PTRACE
# NoNewPrivileges=false

# BPF requires locked memory for maps and ring buffer
LimitMEMLOCK=infinity

# Resource limits
CPUQuota=35%
MemoryMax=512M

# Hardening
ProtectSystem=strict
ReadWritePaths=/var/log/aurora-linux /opt/aurora-linux
PrivateTmp=true

[Install]
WantedBy=multi-user.target
